# services/axon-gui-agent/axon-gui-agent.service
[Unit]
Description=Axon OS GUI Agent (desktop automation) Service
After=dbus.socket

[Service]
Type=simple
ExecStart=/usr/bin/python3 AXON_SERVICES_DIR/axon-gui-agent/gui_agent_service.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
