# services/axon-voice/axon-voice.service
[Unit]
<<<<<<< HEAD
Description=Axon OS AI Voice Service
After=dbus.socket
=======
Description=Axon OS Voice (speech-to-intent) Service
After=graphical-session.target
PartOf=graphical-session.target
>>>>>>> origin/main

[Service]
Type=simple
ExecStart=/usr/bin/python3 AXON_SERVICES_DIR/axon-voice/voice_service.py
Restart=on-failure
RestartSec=3

[Install]
<<<<<<< HEAD
WantedBy=default.target
=======
WantedBy=graphical-session.target
>>>>>>> origin/main
