# services/axon-voice/axon-voice.service
[Unit]
Description=Axon OS Voice (speech-to-intent) Service
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 AXON_SERVICES_DIR/axon-voice/voice_service.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
